# examples/xmrexplorer.conf
[program:xmrexplorer]
command = /srv/xmrexplorer/bin/gunicorn -t 180 -w 4 -b 127.0.0.1:8000 --chdir=/srv/xmrexplorer/ xmrExplorer:app
autostart = true
autorestart = true
